# root/etc/services.d/mailhog/run
#!/usr/bin/execlineb -P

s6-setuidgid abc
exec /usr/local/bin/MailHog
